# .init/sentry.conf
start on runlevel

respawn

exec /app/env/bin/sentry --config=/app/current/sentry_conf.py start 
