# .init/sentry.conf
start on runlevel

respawn

exec sentry --config=/app/sentry/sentry.conf.py start --no-daemon
