# .init/sentry.conf
start on runlevel

respawn

exec /app/env/bin/sentry --config=/app/sentry/sentry.conf.py start 
