# .init/sentry.conf
start on runlevel

respawn

exec sentry start --no-daemon --config=/app/sentry/sentry_conf.py    
